# scripts/crontab
25      2       *       *       *       cd /base/rules/myself && git pull